# templates/default/tickr-upstart.conf.erb
#!upstart
description "Tickr Web Server"

start on startup
stop on shutdown

exec sudo -u root sh -c "<%= @cmd %>"
respawn
